# acpid.service
[Unit]
Description=ACPI Event Daemon
After=syslog.target

[Service]
Type=simple
ExecStart=/usr/sbin/acpid -f

[Install]
WantedBy=multi-user.target
